# raiz/etc/hotplug.d/usb/20-mjpg-streamer
case "$ACTION" in
	add)
	        # start process
		echo "hola mundo" > /root/prendi
		/etc/init.d/mjpg-streamer start
	        ;;
	remove)
	        # stop process
		echo "hola mundo" > /root/apage
		/etc/init.d/mjpg-streamer stop
	        ;;
esac

